# crontab
30 2 * * * root python /crawler.py >> /var/log/cron.log 2>&1
