# crontab
0 0 */15 * * root python /srv/crawler.py > /var/log/cron.log 2>&1
